# deploy/cloudflared-quick.service
[Unit]
Description=Cloudflare quick tunnel for miniapp (writes URL to miniapp_url.txt)
After=network-online.target robokassa-server-denis.service
Wants=network-online.target

[Service]
Type=simple
User=enhel-method
WorkingDirectory=/home/enhel-method
Environment=MINIAPP_URL_FILE=/home/enhel-method/tg-ai-denis-komkov/miniapp_url.txt
Environment=TZ=Europe/Moscow
Environment=PATH=/usr/bin:/bin
ExecStart=/home/enhel-method/tg-ai-denis-komkov/deploy/run_cloudflared_quick.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
